# docker/supervisord.conf
[supervisord]
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid
nodaemon=true
loglevel=debug
user=root

[program:prep]
autostart=false
user=www-data
command=/bin/sh /var/www/docker/launch.sh

[program:php-fpm]
pidfile=/var/run/php-fpm.pid
command=php-fpm
autostart=true
user=root

[program:nginx]
pidfile=/var/run/nginx.pid
command=/usr/sbin/nginx
autostart=true
user=root

[program:worker]
directory=/var/www
autostart=true
user=www-data
command=/usr/local/bin/php artisan queue:work
